# backend/vidwiz/infra/systemctl/vidwiz-transcript-helper.service
[Unit]
Description=Vidwiz Transcript Helper Service

[Service]
EnvironmentFile=/home/adhiraj/projects/vidwiz/.env
ExecStart=/usr/bin/python /home/adhiraj/projects/vidwiz/backend/vidwiz/workers/scripts/transcript-helper.py
Restart=always
Environment=PYTHONUNBUFFERED=1
User=adhiraj
WorkingDirectory=/home/adhiraj/projects/vidwiz/backend

[Install]
WantedBy=multi-user.target